# profile.d/aliases.sh
#!/bin/sh

alias v='vim'
alias l='ls -1'
alias lsa='ls -a'
alias la='ls -1a'
alias ll='ls -hl'
alias lla='ls -alh'

alias g='git'
alias gn='git --no-pager'
alias cdg='cd $(git rev-parse --show-toplevel)'

alias be="bundle exec"

alias d="sdcv --color"
alias hs='fc -l -25'
alias hsa='fc -l 1'
