# profile.d/aliases.sh
#!/bin/sh

alias v='vim'
alias l='ls -1'
alias lf='ls -1F'
alias lsf='ls -F'
alias lsa='ls -a'
alias lsaf='ls -aF'
alias la='ls -1a'
alias laf='ls -1aF'
alias ll='ls -hl'
alias llf='ls -hlF'
alias lla='ls -alh'
alias llaf='ls -alhF'

alias g='git'
alias gn='git --no-pager'
alias cdg='cd $(git rev-parse --show-toplevel)'

alias be="bundle exec"

alias d="sdcv --color"
alias hs='fc -l -25'
alias hsa='fc -l 1'
