# profile.d/linux.sh
if test "$(uname)" != 'Linux' ;then
  return 0
fi

