# profile.d/tty.sh
#!/bin/sh

stty stop undef  # C-S
stty start undef # C-Q
